# .config/systemd/user/clash.service
[Unit]
After = sockets.target

[Service]
ExecStart = %h/notes/linux/clash-wrapper.sh

[Install]
WantedBy = default.target
