# .config/systemd/user/mpd.service
[Unit]
After = sockets.target
After = pipewire-pulse.socket

[Service]
ExecStart = %h/bin/mpd --no-daemon

[Install]
WantedBy = default.target
